# systemd/system/xt_geoip-block_countries.timer
[Unit]
Description=Weekly updates to geoip database from geolite2 data

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
